# stow/chromium-bookmarks/Library/LaunchAgents/com.user.chromium-bookmarks-sync.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.user.chromium-bookmarks-sync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>__HOME__/.local/bin/chromium-bookmarks-sync.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/chromium-bookmarks-sync.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/chromium-bookmarks-sync.log</string>
</dict>
</plist>
