# stow/bin/Library/LaunchAgents/com.user.check-stale-dev-servers.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.check-stale-dev-servers</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>__HOME__/.local/bin/check-stale-dev-servers.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>
	<key>StandardOutPath</key>
	<string>/tmp/check-stale-dev-servers.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/check-stale-dev-servers.log</string>
</dict>
</plist>
