# stow/devonthink/Library/LaunchAgents/com.user.entity-filing.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.user.entity-filing</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>__HOME__/.local/bin/entity-filing.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/entity-filing.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/entity-filing.log</string>
	<key>StartInterval</key>
	<integer>1800</integer>
</dict>
</plist>
