# stow/happy/Library/LaunchAgents/com.user.happy-daemon.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.happy-daemon</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>exec __HOME__/.local/share/mise/shims/happy daemon start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/happy-daemon.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/happy-daemon.log</string>
</dict>
</plist>
